# assets/hypr-session-daemon.service
[Unit]
Description=hypr-session auto-save daemon
PartOf=graphical-session.target
After=graphical-session.target

[Service]
Type=simple
ExecStart=%h/.local/bin/hypr-session daemon
Restart=on-failure
RestartSec=5

[Install]
WantedBy=graphical-session.target
